# deploy/evo-api.service
[Unit]
Description=AUTO-EVO-AI V0.1 API Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/my-evo-ai
EnvironmentFile=/home/ubuntu/my-evo-ai/.env
ExecStart=/home/ubuntu/my-evo-ai/venv/bin/python3 -m uvicorn api_server:app --host 0.0.0.0 --port 8765 --log-level info
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/ubuntu/evo.log
StandardError=append:/home/ubuntu/evo.log

# ○ 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
